# skitai/install-data/etc/init/skitaid.conf
description "Skitai App Engine"
author "Hans Roh <hansroh@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

chdir /data/skitai/skitai/bin

exec /usr/local/bin/skitaid.py

respawn
