# skitai/skitaid/etc/init.d/skitaid
#!/bin/sh
# chkconfig: 123456 90 10
# TTS Server for Speech Synthesis
#
workdir=/home/skitai/skitai
 
start() {
    skitaid.py &
    echo "Skitaid App Engine Started."
}
 
stop() {
    pid=`ps -ef | grep '[s]kitaid.py' | awk '{ print $2 }'`
    #echo $pid
    kill $pid
    sleep 2
    echo "Skitai App Engine Killed."
}
 
case "$1" in
  start)
    start
    ;;
  stop)
    stop   
    ;;
  restart)
    stop
    start
    ;;
  *)
    echo "Usage: /etc/init.d/skitaid {start|stop|restart}"
    exit 1
esac
exit 0
