# examples/resources/init-skitaid.conf
description "Skitai App Engine"
author "Hans Roh <hansroh@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn
exec /home/ubuntu/skitai/exapmples/app.py -d
